# Cti.ini
[program:cti]
directory = /data/www/cti_esl/Release
command = dotnet Cti.dll
autostart = true
stdout_logfile = /data/www/cti_esl/Release/runtime/runtime.log
redirect_stderr = true